# script/shell/systermrun/source.sh
#!/bin/bash
#ScriptName: source

### BEGIN INIT INFO
# Provides:          Bec
# Required-Start:    None
# Required-Stop:     None
# Default-Start:     None
# Default-Stop:      None
# Short-Description: source /etc/profile
# Description:       source /etc/profile
### END INIT INFO

source /etc/profile
